# templates/etc/systemd/system/gitlab-runner.service.d/kill.conf.j2
[Service]
TimeoutStopSec={{ runner_stop_timeout }}
KillSignal=SIGQUIT
